# src/guest-docker.path
[Unit]
[Path]
PathExistsGlob=/dev/disk/by-path/*guest-docker*
[Install]
WantedBy=default.target
